# bootstrap/systemd/firstboot.service
[Unit]
Description=Run first boot cleanup once
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/firstboot.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
